# ansible/roles/api/files/sd_api.service
[Unit]
Description=Stable Diffusion API

[Service]
ExecStart=/home/ubuntu/sd_api.sh
Restart=on-failure
User=ubuntu
